# deploy/bootc/rhoim-model-pull.service
[Unit]
Description=RHOIM Model Pull
After=network-online.target
Before=rhoim-vllm.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/rhoim
ExecStart=/usr/bin/bash -lc '/opt/rhoim/bin/model_pull.sh'
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
