# deploy/bootc/rhoim-gateway.service
[Unit]
Description=RHOIM Gateway
After=network-online.target rhoim-vllm.service

[Service]
EnvironmentFile=-/etc/sysconfig/rhoim
Environment=VLLM_ENDPOINT=http://127.0.0.1:8000
Environment=BASE_PATH=/
Environment=ENABLE_RHOAI_ALIAS=true
ExecStart=/opt/rhoim/venv/bin/python3 -m uvicorn app.main:app --host 0.0.0.0 --port 8080 --app-dir /opt/rhoim/gateway
WorkingDirectory=/opt/rhoim/gateway
Restart=always

[Install]
WantedBy=multi-user.target
